# .config/systemd/user/init-user-env.service
[Unit]
Description=User Environment Initialization
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=%h/.local/share/bin/init-user-env.sh