# .config/systemd/user/nix-collect-garbage.service
[Unit]
Description=Nix Collect Garbage
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=nix-collect-garbage
