# .config/systemd/user/podman-cleanup.service
[Unit]
Description=Podman System Cleanup
ConditionACPower=true

[Service]
Type=oneshot
ExecStart=podman system prune -a --volumes -f
